# etc/systemd/system/b.service
[Unit]
Description=do airconsystem

[Service]
ExecStart=/usr/bin/python3 /home/pi/systemd/getTimeline.py
